# pi/deploy/hubd.service
[Unit]
# Generic Linux unit — runs hubd as the always-on classroom dashboard/HTTP
# chassis on any systemd host (a Pi is just one). hubd is NOT the transport:
# that's zenohd + the ws-adapter, separate services (deploy/install.sh).
# hubd serves the dashboard (:80), /fleet, and device-served Wi-Fi setup
# (/wifi/*, driving NetworkManager) — the latter folded in when provisiond and
# the BLE stack were deleted (2026-07-09).
Description=hub — classroom dashboard + Wi-Fi setup (transport is zenohd + ws-adapter)
# NetworkManager (not just network-online): /wifi/scan|connect shell out to nmcli.
After=network-online.target NetworkManager.service
Wants=network-online.target NetworkManager.service

[Service]
Type=simple
WorkingDirectory=/opt/hub
# Dashboard on :80 so the classroom URL is just http://hub.local (avahi/mDNS,
# hostname `hub`; 10.42.0.1 is the fallback). Code default stays :8000 so
# unprivileged dev runs (`cargo run --bin hubd`) still bind. hubd holds no
# transport session, so it needs no credentials.
Environment=HUB_HTTP=0.0.0.0:80
ExecStart=/opt/hub/hubd
Restart=always
RestartSec=2
# Root, not the old DynamicUser: hubd now drives NetworkManager (nmcli) for the
# device-served Wi-Fi panel, which a dynamic UID can't do (no D-Bus/polkit
# identity). This is the privilege the deleted provisiond used to hold; folding
# provisioning into hubd folds that requirement in too. The blast radius is a
# classroom-LAN HTTP server that can change the hub's Wi-Fi — which is exactly
# the feature, so the grant is the capability, not a leak. :80 needs no separate
# CAP_NET_BIND_SERVICE as root.
User=root

[Install]
WantedBy=multi-user.target
