# pi/deploy/provisiond.service
[Unit]
Description=hub Wi-Fi provisioning over Improv BLE
# BlueZ for GATT/advertising, NetworkManager for the actual join.
After=bluetooth.target NetworkManager.service
Wants=bluetooth.target

[Service]
ExecStart=/opt/hub/provisiond
Restart=always
RestartSec=2
# Root, not DynamicUser (unlike hubd): the local GATT app over BlueZ's D-Bus,
# the legacy advertisement via btmgmt's kernel mgmt socket (CAP_NET_ADMIN), and
# driving nmcli all need privilege. The non-root paths are brittle across BlueZ
# versions. (btmgmt: the Pi 4 controller can't do bluetoothd's ext-advertising —
# see provisiond.rs start_legacy_advert.)
User=root

[Install]
WantedBy=multi-user.target
